# deploy/ground-up-wall.service
[Unit]
Description=ground-up-wall photowall
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=groundupwall
Group=groundupwall
WorkingDirectory=/opt/ground-up-wall
EnvironmentFile=/opt/ground-up-wall/.env
ExecStart=/usr/local/bin/deno run -A /opt/ground-up-wall/prod.ts
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/ground-up-wall /opt/ground-up-wall

[Install]
WantedBy=multi-user.target
